# supervisord/filer.ini
[program:filer]
process_name = filer
command = /sites/filer/filer-backend
directory = /sites/filer/
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/filer/logs/stdout.log
environment = HOME="/root"