# supervisord/imaginary.ini
[program:imaginary]
process_name = imaginary
command = /sites/imaginary/start.sh
directory = /sites/imaginary/
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/imaginary/logs/imaginary.log