# supervisord/md-fs.ini
[program:mdb-fs]
process_name = mdb-fs
command = /sites/mdb-fs/mdb-fs
directory = /sites/mdb-fs/
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/mdb-fs/logs/mdb-fs.log